# home/dot_config/autostart/1password.desktop
[Desktop Entry]
Exec=/usr/bin/1password --silent
Name=1password
Type=Application
Icon=1password
Terminal=false
StartupNotify=false
