# home/dot_config/autostart/todoist.desktop
[Desktop Entry]
Exec=/var/lib/snapd/snap/bin/todoist %U
Name=todoist
Type=Application
Icon=todoist
Terminal=false
Type=Application
